# docker/supervisor/conf.d/erp-consumer.conf
[program:erp-consume]
command=php /var/www/artisan erp:consume
autostart=true
autorestart=true
user=root
redirect_stderr=true
stdout_logfile=/var/www/storage/logs/erp-consume.log
stopwaitsecs=3600
